# cmd/imaginator/example-manifests/dominator/files/etc/systemd/system/filegen-server.service.d/env.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/filegen-server /var/lib/filegen-server/computed-files/root
User=filegen-server
