# init.d/mdbd.Debian-7.0
#! /bin/bash --posix

### BEGIN INIT INFO
# Provides:		mdbd
# Required-Start:	$local_fs $network $syslog
# Required-Stop:	$local_fs $network $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description:	MDB query daemon
### END INIT INFO

set -e

# /etc/init.d/mdbd: start and stop the MDB query daemon

. /lib/lsb/init-functions

umask 022

DAEMON='/usr/local/sbin/mdbd'
MDB_FILE='/var/lib/Dominator/mdb'
PIDFILE='/var/run/mdbd.pid'
SOURCES_FILE='/var/lib/Dominator/mdb.sources.list'
USERNAME='mdbd'

[ -f /etc/default/mdbd ] && . /etc/default/mdbd

test -x "$DAEMON" || exit 0

export PATH="${PATH:+$PATH:}/usr/local/bin:/usr/local/sbin:/usr/sbin:/sbin"

case "$1" in
  start)
	log_daemon_msg "Starting MDB query daemon" "mdbd" || true
	start-stop-daemon --start --quiet --pidfile "$PIDFILE" \
			  --background --exec "$DAEMON" --chuid "$USERNAME" \
			  --make-pidfile -- \
			  -mdbFile="$MDB_FILE" -sourcesFile="$SOURCES_FILE"
	;;
  stop)
	log_daemon_msg "Stopping MDB query daemon" "mdbd" || true
	kill -TERM $(cat "$PIDFILE")
	rm -f "$PIDFILE"
	;;

  reload|force-reload)
	$0 stop
	$0 start
	;;

  restart)
	$0 stop
	$0 start
	;;

  *)
	log_action_msg "Usage: /etc/init.d/mdbd {start|stop|reload|force-reload|restart}" || true
	exit 1
esac

exit 0
